# root/etc/services.d/beanstalkd/run
#!/usr/bin/execlineb -P

beanstalkd -p 11300 -u nobody
